# bruno/gfhd100/skel/etc/inittab
# /etc/inittab
#
# This inittab is a basic inittab sample for sysvinit, which mimics
# Buildroot's default inittab for Busybox.
id:1:initdefault:

run1:1:wait:/etc/init.d/rcS start
S0:1:respawn:/sbin/rungetty

# Stuff to do for the 3-finger salute
rebo::ctrlaltdel:/sbin/reboot

# Stuff to do before reboot/shutdown
run6:6:wait:/etc/init.d/rcS stop

# Actually reboot/shutdown
rebf:6:wait:/sbin/reboot -f
shut:0:wait:/sbin/halt -p -f
